# slave/root/etc/services.d/sshd/run
#!/usr/bin/with-contenv bash

/usr/sbin/sshd -D

## For debugging
#/usr/sbin/sshd -D -ddd -E /var/log/ssh.log
